# systemd/screensaver-stop.service
[Unit]
Description=Stop screensaver

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop screensaver.service
