# systemd/screensaver-stop.timer
[Unit]
Description=Stop screensaver at 05:25

[Timer]
OnCalendar=*-*-* 05:25:00
AccuracySec=1s

[Install]
WantedBy=timers.target
