# systemd/screensaver-start.timer
[Unit]
Description=Start screensaver at 09:05

[Timer]
OnCalendar=*-*-* 09:05:00
AccuracySec=1s

[Install]
WantedBy=timers.target
