# remote/systemd/companion-remote.service
[Unit]
Description=companion remote (スマホ専用リモコン PWA + 配信サーバ, 127.0.0.1)
Documentation=file://%h/companion/remote/docs/STATUS.md
Documentation=file://%h/companion/workspace/redesign/remote-design.md

[Service]
Type=simple
WorkingDirectory=%h/companion/remote
# .env は任意(- 付きで不在でも起動)。REMOTE_PORT を渡す。
EnvironmentFile=-%h/companion/remote/.env
# (F-2 発話撤去[2026-06-14]に伴い Environment=DISPLAY=:0 を削除。remote プロセスは
#  描画も audio も持たず DISPLAY 不要 = say.sh paplay 経路専用だったため。STATUS.md Done 参照)
ExecStart=%h/companion/remote/server/app.py
# 常駐 HTTP サーバなのでクラッシュ時は再起動(stateless、回復に state 引きは不要)。
# always: 外部からの SIGTERM は systemd では正常終了扱いで on-failure では復帰しない
# (2026-07-21 pkill 巻き込み停止)。意図停止 (systemctl stop) は always でも再起動されない。
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
